# operations/hetzner-production-live-sync/systemd/sellton-postgres-backup-standby.service
[Unit]
Description=Sellton production PostgreSQL post-cutover backup checkpoint
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
User=root
Group=root
ExecStart=/opt/sellton/live-sync/postgres-standby/postgres-standby-backup.sh
