# operations/hetzner-production-live-sync/systemd/sellton-storage-standby.service
[Unit]
Description=Sellton production Supabase Storage Hetzner to cloud standby sync
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service
Conflicts=sellton-storage-live-sync.service

[Service]
Type=oneshot
User=root
Group=root
ExecStart=/opt/sellton/live-sync/storage-standby/storage-standby-sync-run.sh sync
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=6
